# 99seal-tpm/unseal.service.in
[Unit]
Description=Unseal TPM-sealed files
DefaultDependencies=no
Before=BEFORE
After=tpm2.target systemd-pcrphase.service systemd-pcrphase-initrd.service
Requires=tpm2.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/unseal.sh

[Install]
RequiredBy=REQUIRED_BY
WantedBy=WANTED_BY
